# docker/supervisor/queue.conf
[program:laravel-worker]
command=php /var/www/html/artisan queue:work
autostart=true
autorestart=true
